# src/.config/systemd/user/mail-sync.service
[Unit]
Description=Sync mail

[Service]
ExecStart=/bin/sh -c "mbsync -a && notmuch new"
