# src/.config/systemd/user/mail-sync.service
[Unit]
Description=Sync mail

[Service]
ExecStart=mbsync -a
