# src/.config/systemd/user/mail-sync.timer
[Unit]
Description=Sync the mailbox every 10 minutes

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
