# src/.config/systemd/user/votecoind.service
[Unit]
Description=VoteCoin daemon

[Service]
ExecStart=%h/.local/bin/votecoind -metricsui 0 -daemon

[Install]
WantedBy=default.target

